# dev-db/mha-manager/files/mha-manager-initd
#!/sbin/runscript

extra_commands="checkconfig"

depend() {
        use net
}

start() {
        checkconfig || return 1

	einfo "Starting with ${mha_options}"

        ebegin "Starting MHA manager"
        start-stop-daemon \
          --start \
          --quiet \
          -1 ${MHA_LOG_STDOUT}  \
          -2 ${MHA_LOG_STDERR} \
          --make-pidfile \
          --pidfile /run/mha-manager.pid \
          --background \
	  --wait 10000 \
          --exec /usr/bin/masterha_manager -- ${mha_options}
        eend $?
}

stop() {
        ebegin "Stopping MHA manager"
        start-stop-daemon \
          --stop \
          --quiet \
          --pidfile=/run/mha-manager.pid
        eend $?
}

checkconfig() {
        mha_options="--conf=${MHA_APP_CONFIGFILE} --global_conf=/etc/MHA/masterha_default.cnf"

	if [ -z "${MHA_LAST_FAILOVER_MINUTES}" ]; then
	  mha_options="${mha_options} --ignore_last_failover"
	else
	  mha_options="${mha_options} --last_failover_minute=${MHA_LAST_FAILOVER_MINUTES}"
	fi

	if [ "${MHA_IGNORE_FAIL_ON_START}" == "yes" ]; then
	  mha_options="${mha_options} --ignore_fail_on_start"
	fi

        ebegin "Checking configuration sanity"

        if [ ! -r "${MHA_APP_CONFIGFILE}" ]; then
          eerror "MHA_APP_CONFIGFILE is set to an unreadable file!"
          eerror "Please fix your /etc/conf.d/mha-manager file."
          eend 1
          return 1
        fi

        #/usr/bin/masterha_check_ssh --conf=${MHA_APP_CONFIGFILE} --global_conf=/etc/MHA/masterha_default.cnf || (
        /usr/bin/masterha_check_ssh ${mha_options} || (
          eerror "SSH connection check failed! Please make sure that this host can connect"
          eerror "Please make sure that this host can reach your database servers over ssh."
          eend 1
          return 1
        )

        /usr/bin/masterha_check_repl ${mha_options} || (
        #/usr/bin/masterha_check_repl --conf=${MHA_APP_CONFIGFILE} --global_conf=/etc/MHA/masterha_default.cnf || (
          eerror "Your replication setup is not working."
          eerror "Please fix it first."
          eend 1
          return 1
        )

        eend 0
}
